# energy-monitor.service
[Unit]
Description=Energy Monitor - FastAPI backend
After=network.target

[Service]
Type=simple
# User=pi is default for Raspberry Pi OS; change if needed (e.g. User=%i for your login)
User=pi
Group=pi
WorkingDirectory=REPLACE_INSTALL_DIR
Environment="PYTHONPATH=REPLACE_INSTALL_DIR/backend"
Environment="ENERGY_MONITOR_DB_PATH=REPLACE_INSTALL_DIR/data"
Environment="PYTHONUNBUFFERED=1"
ExecStart=REPLACE_INSTALL_DIR/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
StandardOutput=journal
StandardError=journal
SyslogIdentifier=energy-monitor
Restart=unless-stopped
RestartSec=5

[Install]
WantedBy=multi-user.target
